# systemd/user/battery-alert.service
[Unit]
Description=Battery alert

[Service]
ExecStart=%h/.local/bin/user/battery_notification

[Install]
WantedBy=multi-user.target
